# platform/deploy/supervisord.conf
; =============================================================================
; JustSoc Platform 容器进程管理
;   - platform : Next.js standalone (:3000)
;   - bridge   : Claude Code Bridge (127.0.0.1:4317, 仅容器内)
; 设置页「重启 Bridge」-> supervisorctl restart bridge
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
; 让以 node 用户运行的 Next.js 进程也能调用 supervisorctl restart bridge
chown=node:node

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:platform]
command=node server.js
directory=/app
user=node
; supervisord 以 root 运行、用 user=node 派生子进程时不会重置 HOME（默认仍是 /root，
; node 无权访问），必须显式设 HOME=/home/node，否则 Claude Code CLI 等读写 ~/.claude 会失败。
; PORT 也必须按进程分别指定，否则两个进程会争用同一端口。
environment=HOME="/home/node",USER="node",LOGNAME="node",PORT="3000",HOSTNAME="0.0.0.0",NODE_ENV="production"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:bridge]
command=node dist/server.js
directory=/app/claude-code-bridge
user=node
; HOME 必须显式设为 /home/node：bridge 会 spawn Claude Code CLI，CLI 需读写 ~/.claude，
; 否则继承 root 的 HOME=/root（无权访问）会导致 CLI 启动即静默退出、任务空回复。
environment=HOME="/home/node",USER="node",LOGNAME="node",PORT="4317",NODE_ENV="production"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
